# etc/profile.d/misc.sh
export PATH=$PATH:~/.local/bin
export RXVT_SOCKET=/tmp/urxvtd-socket
