# etc/profile.d/misc.sh
export RXVT_SOCKET=/tmp/urxvtd-socket
export GIT_DISCOVERY_ACROSS_FILESYSTEM=1
