# etc/profile.d/xdg_cache_home.sh
# To spare SSD
export XDG_CACHE_HOME="/tmp/$USER/.cache"
# Exclude sbt-coursier plugin
export COURSIER_CACHE="/home/$USER/.coursier/cache/v1"
