# etc/profile.d/xdg_cache_home.sh
# To spare SSD
export XDG_CACHE_HOME="/tmp/$USER/.cache"
