# etc/local.d/vim_view.start
#!/bin/sh

vim_view=/tmp/vim_view
mkdir $vim_view
chmod a+rw $vim_view
