# etc/udev/rules.d/power.rules
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_TYPE}=="Mains",RUN+="/usr/local/bin/vicious-force batwidget"
